# files/etc/uci-defaults/99-custom
#!/bin/sh

uci set network.lan.ipaddr='192.168.255.254'
uci set network.lan.netmask='255.255.255.0'
uci commit network

/etc/init.d/network restart

exit 0